# deploy/systemd/addref.service.example
[Unit]
Description=AddRef web application
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/addref
ExecStart=/usr/bin/python3 /opt/addref/server.py
Restart=always
RestartSec=3
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
